# vm/files/etc/init.d/metadata
#!/sbin/openrc-run

description="metadata probing"

depend() {
    after network-wait
    provide metadata
}

start() {
    usr/bin/metadata aws
    cp /run/config/conf.d/* /etc/conf.d
    return 0
}