# vm/files/etc/init.d/minecraft
#!/sbin/openrc-run

description="Managing the Minecraft server"

depend() {
	after metadata
}

start() {
	/usr/bin/minecraft &
	/usr/bin/minecraft-sleep &
}