# vm/files/etc/init.d/status-reporter
#!/sbin/openrc-run

description="Report status to the GraphQL API"
output_log="/var/log/status-reporter.log"
error_log="/var/log/status-reporter.err"
command_background="yes"
pidfile="/tmp/status-reporter.pid"

depend() {
	after minecraft
}

stop() {
    /usr/bin/curl \
  		-X POST \
  		-H "Content-Type: application/json" \
		-H "Authorization: token "${GRAPHQL_TOKEN}"" \
  		--data "\"query\": { mutation(destroyServer(id: "$(curl http://169.254.169.254/v1/instanceid)")) }" \
  		https://new.walrus.gg # move that to the script
}