# vm/files/etc/conf.d/minecraft
#!/sbin/openrc-run

description="Managing the Minecraft server"
user=${user:-root}
group=${group:-root}

command="/usr/bin/minecraft"

depend() {
	need net
	after net
}